# scripts/linux/common_tpl.service
## 注意此文件的文件名保持为common_tpl.service，不要调整成其他的。如果要调整修改install-linux.sh里引用部分
[Unit]
Description={{APP_DESCRIPTION}}
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
ExecStart={{APP_NAME}}
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

# 工作目录
WorkingDirectory={{INSTALL_DIR}}

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths={{INSTALL_DIR}}

[Install]
WantedBy=multi-user.target
